# photobooth.service
[Unit]
Description=TVS Photobooth (payment daemon + Pygame app)
After=graphical.target sound.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=tvs
Group=tvs
WorkingDirectory=/home/tvs/TVS_PHOTOBOOTH_V2
Environment=PYTHONUNBUFFERED=1
Environment=HOME=/home/tvs
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/tvs/.Xauthority
ExecStart=/bin/bash /home/tvs/TVS_PHOTOBOOTH_V2/run.sh
Restart=always
RestartSec=5
# run.sh already forwards SIGTERM → cleanup → kill children
KillSignal=SIGTERM
TimeoutStopSec=15
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical.target
